# etc/systemd/system/apt-update.timer
[Unit]
Description=Run apt update nightly

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
