# etc/systemd/system/apt-update.service
[Unit]
Description=Nightly apt update

[Service]
Type=oneshot
ExecStart=/usr/bin/apt update

[Install]
WantedBy=multi-user.target
